# server/com.swingcam.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.swingcam.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ryancase/dev/garage_swing_detector/server/venv/bin/python</string>
        <string>/Users/ryancase/dev/garage_swing_detector/server/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ryancase/dev/garage_swing_detector/server</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <!-- Restart if it crashes, but respect manual stop -->
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ryancase/dev/garage_swing_detector/server/logs/swingcam.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ryancase/dev/garage_swing_detector/server/logs/swingcam.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
